# contrib/initscripts/bsd/bitbean
 
#!/bin/sh

# PROVIDE: bitbean
# REQUIRE: LOGIN
# KEYWORD: shutdown

#
# Add the following lines to /etc/rc.conf.local or /etc/rc.conf
# to enable this service:
#
# bitbean_enable (bool): Set to NO by default. Set this variable to YES if you want to enable bitbean service.
# bitbean_config (path): Set to /usr/local/etc/bitbean.conf by default.
# bitbean_user:  The user account bitbean daemon runs as. It uses 'root' user by default.
# bitbean_group: The group account bitbean daemon runs as. It uses 'wheel' group by default.
# bitbean_datadir (str): Default to "/var/db/bitbean". Base data directory.

. /etc/rc.subr

name=bitbean
rcvar=bitbean_enable

: ${bitbean_enable:=NO}
: ${bitbean_config=/usr/local/etc/bitbean.conf}
: ${bitbean_datadir=/var/db/bitbean}
: ${bitbean_user="root"}
: ${bitbean_group="wheel"}

required_files=${bitbean_config}
command=/usr/local/bin/bitbeand
bitbean_chdir=${bitbean_datadir}
pidfile="${bitbean_datadir}/bitbeand.pid"
stop_cmd=bitbean_stop
command_args="-conf=${bitbean_config} -datadir=${bitbean_datadir} -daemon -pid=${pidfile}"
start_precmd="${name}_prestart"

bitbean_create_datadir()
{
	echo "Creating data directory"
	eval mkdir -p ${bitbean_datadir}
	[ $? -eq 0 ] && chown -R ${bitbean_user}:${bitbean_group} ${bitbean_datadir}
}

bitbean_prestart()
{
	if [ ! -d "${bitbean_datadir}/." ]; then
		bitbean_create_datadir || return 1
	fi
}

bitbean_requirepidfile()
{
	if [ ! "0`check_pidfile ${pidfile} ${command}`" -gt 1 ]; then
		echo "${name} not running? (check $pidfile)."
		exit 1
	fi
}

bitbean_stop()
{
    bitbean_requirepidfile

	echo "Stopping ${name}."
	eval ${command} -conf=${bitbean_config} -datadir=${bitbean_datadir} stop
	wait_for_pids ${rc_pid}
}

load_rc_config $name
run_rc_command "$1"
